# deploy/systemd/grookai-tcgcsv-historical-backfill.service
[Unit]
Description=Grookai TCGCSV Historical Archive Backfill Agent
Wants=network-online.target
After=network-online.target grookai-mee-reference-refresh.service grookai-mee-nightly.service grookai-tcgcsv-warehouse.service

[Service]
Type=oneshot
User=grookai
WorkingDirectory=/opt/grookai_vault_mee_nightly
EnvironmentFile=/etc/grookai/mee-nightly.env
Environment=TCGCSV_HISTORICAL_BACKFILL_START_DATE=2024-02-17
Environment=TCGCSV_HISTORICAL_BACKFILL_STATE_DIR=/var/lib/grookai
Environment=TCGCSV_HISTORICAL_BACKFILL_OUT_DIR=docs/audits/market_evidence_engine_v1/tcgcsv_full_source_warehouse_v1
Environment=TCGCSV_HISTORICAL_BACKFILL_PAUSE_START_HHMM=0050
Environment=TCGCSV_HISTORICAL_BACKFILL_PAUSE_END_HHMM=1030
Environment=TCGCSV_HISTORICAL_BACKFILL_CLEANUP_EXTRACTED=1
Environment=TCGCSV_HISTORICAL_BACKFILL_MAX_DAYS=1
Environment=TCGCSV_HISTORICAL_BACKFILL_MAX_RUNTIME_SECONDS=7200
Environment=TCGCSV_HISTORICAL_BACKFILL_LOAD_GUARD_ENABLED=1
ExecStart=/usr/bin/env bash scripts/workers/tcgcsv_historical_backfill_agent_v1.sh
Restart=no
TimeoutStartSec=2h
TimeoutStopSec=2m
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
OOMScoreAdjust=500
