# deploy/systemd/grookai-tcgcsv-warehouse.service
[Unit]
Description=Grookai TCGCSV Full Source Warehouse Current Sync
Wants=network-online.target
After=network-online.target grookai-mee-nightly.service

[Service]
Type=oneshot
User=grookai
WorkingDirectory=/opt/grookai_vault_mee_nightly
EnvironmentFile=/etc/grookai/mee-nightly.env
Environment=TCGCSV_FULL_SYNC_REQUEST_CEILING=10000
Environment=TCGCSV_REQUEST_DELAY_MS=100
Environment=TCGCSV_DIMENSION_BATCH_SIZE=1000
Environment=TCGCSV_ARTIFACT_BATCH_SIZE=1000
Environment=TCGCSV_PRICE_OBSERVATION_BATCH_SIZE=500
ExecStart=/usr/bin/flock -n /tmp/grookai-tcgcsv-warehouse.lock /usr/bin/env bash -lc 'node scripts/workers/tcgcsv_full_source_warehouse_worker_v1.mjs --mode=current --apply --out-dir=docs/audits/market_evidence_engine_v1/tcgcsv_full_source_warehouse_v1'
TimeoutStartSec=4h
Nice=8
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
