# deploy/cron/grookai-mee-nightly.cron
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

15 3 * * * cd /opt/grookai_vault && set -a && source /etc/grookai/mee-nightly.env && set +a && node scripts/workers/mee_nightly_droplet_worker_v1.mjs --run --call-ceiling=${MEE_NIGHTLY_CALL_CEILING:-4000} >> /var/log/grookai-mee-nightly.log 2>&1
